# systemd/user/subniri-bar.service
[Unit]
Description=Subniri bar service
PartOf=subniri.target
After=graphical-session.target

[Service]
Type=simple
Environment="GDK_BACKEND=wayland"
ExecStart=polarbar
Restart=on-failure
RestartSec=1s

[Install]
WantedBy=subniri.target
