# systemd/user-dev/subniri-dev-launcher.service
[Unit]
Description=Subniri development launcher service
PartOf=subniri-dev.target
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/projects/subniri/target/debug/avalaunch
WorkingDirectory=%h/projects/subniri
Restart=on-failure
RestartSec=1s

[Install]
WantedBy=subniri-dev.target
